# supervisord.conf
; supervisord-конфиг для max-telegram-bridge.
; Внутри одного контейнера:
;   - xvfb   :99                  (виртуальный X-сервер)
;   - x11vnc  -display :99         (VNC-стрим экрана на 5900)
;   - novnc   websockify 6080     (HTML5 noVNC клиент на /vnc)
;   - api    uvicorn (FastAPI)    (8000)
;   - bot    aiogram              (Telegram)
;   - watcher Playwright + chromium
; Логи шлются в /dev/stdout — критично для managed-хостингов (Railway / Render).

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; --- X-сервер (запускается первым) ---
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x800x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- VNC-стрим экрана :99 ---
[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbauth /etc/x11vnc.passwd -wait 5 -rfbport 5900
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- noVNC: HTML5-обёртка над VNC на порту 6080 ---
[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true
priority=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- API (FastAPI) ---
[program:api]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --proxy-headers
directory=/app/api
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- Telegram-бот ---
[program:bot]
command=python run.py
directory=/app/bot
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- Watcher (Playwright + Chromium) ---
[program:watcher]
command=python run.py
directory=/app/watcher
autostart=true
autorestart=true
priority=30
environment=CHROMIUM_PATH="/usr/bin/chromium",DISPLAY=":99",WATCHER_HEADFUL_DEFAULT="0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0